# scripts/gapapp
#!/sbin/openrc-run

depend () {
	need net
}

command=/usr/bin/gunicorn
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-p ${pidfile} -c /home/mag/GapApp/gap/gunicorn.conf.py --chdir=/home/mag/GapApp/gap/ gap.wsgi"
command_background=true
command_user="mag:mag"

